# infra/dockteacher-worker.service
[Unit]
Description=DOCKTEACHER video transcoding worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/dockteacher-web/worker
EnvironmentFile=/home/ubuntu/dockteacher-web/worker/.env
ExecStart=/usr/bin/node /home/ubuntu/dockteacher-web/worker/index.js
Restart=always
RestartSec=5
# SIGTERM 드레인: 진행 중인 트랜스코딩(최대 수십 분)이 끝날 때까지 강제 kill하지 않는다
TimeoutStopSec=7200

[Install]
WantedBy=multi-user.target
